# ops/systemd/immich-library-scan.service
[Unit]
Description=Trigger an Immich external library scan
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/opt/image/.env
ExecStart=/opt/image/scripts/immich-rescan.sh

